# elements/osism-ipa/init-scripts/systemd/osism-ipa-config.service
[Unit]
Description=OSISM IPA Config Generation

[Service]
ExecStart=/usr/local/sbin/osism-ipa-configure
Type=oneshot
RemainAfterExit=yes
TimeoutSec=30

[Install]
WantedBy=multi-user.target
