# elements/ansible-init/init-scripts/systemd/osism-configdrive.mount
[Unit]
Description=Mount OSISM Configuration Drive
StopWhenUnneeded=yes

[Mount]
What=/dev/disk/by-label/config-2
Where=/mnt/config
